# packaging/linux/systemd/alertpriority-agent.service
[Unit]
After=network.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/alertpriority-agent --logfile /var/log/alertpriority-agent/apagent.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
